# docs/agents/hermes/anchor-timer/fluncle-anchor.service
[Unit]
Description=fluncle-anchor — rave-02 host trigger for the catalogue Spotify-anchor sweep
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/anchor-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service
# If this unit dies before writing its own /status marker (OOM, missing binary, a crash
# before the sweep catch runs), the healthcheck prober never sees a failure. This catch-all
# posts a minimal Discord alert; %n is this unit's full name.
OnFailure=fluncle-sweep-failure@%n.service

[Service]
Type=oneshot
# Root drives the Docker daemon; `docker exec -u hermes` runs the WORK as the unprivileged
# hermes user. The .sh/.ts pair is BAKED at /opt/hermes-scripts/ (auto-updates from main via
# pin-watch); this host unit only TRIGGERS it, and the sweep self-writes its /status marker
# (cron-output.sh). Secrets (the agent token + APIFY_API_TOKEN) come from the shared op-injected
# ${HOME}/.fluncle-secrets.env the wrapper sources.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home -e FLUNCLE_ANCHOR_BATCH=250 hermes bash /opt/hermes-scripts/anchor-sweep.sh
# One tick drains 250 rows (`FLUNCLE_ANCHOR_BATCH`, the `list_track_work?limit` cap) — raised from the
# code default of 15, safe now that the FREE rungs (Deezer ISRC-recovery + ListenBrainz + Spotify
# search) resolve most rows at zero Apify cost, the reason the batch was pinned small. At ~3.7s/row
# (Spotify ≥2s search pacing) a full tick is ~15 min, so the timeout is 30 min; it still sits under the
# 1h cadence, so an over-long tick just skips the next trigger and self-heals. Once the crawl backlog is
# dry the sweep no-ops on an empty queue — lower the batch back toward 15 for steady-state if desired.
TimeoutStartSec=1800
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
