# docs/agents/hermes/capture-timer/fluncle-capture.service
[Unit]
Description=fluncle-capture — rave-02 host trigger for the full-song capture sweep (decoupled from the latency-sensitive gateway cron runner)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/capture-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
# The host unit runs as root (to drive the Docker daemon); `docker exec -u hermes` runs the
# in-container capture WORK as the UNPRIVILEGED hermes user — yt-dlp handles untrusted network
# content, so it must not run as container-root (hermes owns the 10000 secrets file the .sh
# sources). The .sh/.ts pair (+ the pinned yt-dlp) is BAKED into the image at
# /opt/hermes-scripts/ (Unit A/D — auto-updates from main via pin-watch, no /opt/data/scripts
# copy); this host unit only TRIGGERS it. A per-tick yt-dlp fetch through a residential proxy
# has an UNBOUNDED tail (60s search + up to a 180s download, times BATCH_CAP), so it must NOT
# ride the one serial Hermes cron runner (a ~300s global budget) where a backfill-drain
# tick would starve the latency-sensitive 5-min enrich/context/note sweeps — the same
# reason the healthcheck prober and the pin-watch rebuild are host timers.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/capture-sweep.sh
# A worst-case tick (BATCH_CAP findings, each up to a 180s download) can run minutes; give
# it headroom well past the gateway's ~300s global budget.
TimeoutStartSec=600
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
