# docs/agents/hermes/live-timer/fluncle-live.service
[Unit]
Description=fluncle-live — rave-02 host trigger for the Twitch live-set poller (migrated off the Hermes gateway cron runner)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/live-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
# Root drives the Docker daemon; `docker exec -u hermes` runs the WORK as the unprivileged
# hermes user (hermes owns the 0600 secrets file the .sh sources — the Twitch client id/secret).
# The .sh/.ts pair is BAKED at /opt/hermes-scripts/ (Unit A — auto-updates from main via
# pin-watch); this host unit only TRIGGERS it, and the sweep self-writes its /status marker
# (cron-output.sh) so a future prober entry stays honest.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/fluncle-live.sh
# A single Twitch Helix poll (token mint/reuse + one GET + one POST); a clean tick is seconds.
TimeoutStartSec=60
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
