# docs/agents/hermes/artist-follow-timer/fluncle-artist-follow.timer
[Unit]
Description=fluncle-artist-follow timer — trigger the artist auto-follow sweep every 6h on the rave-02 host (not the Hermes gateway)

[Timer]
# First tick shortly after boot (staggered at 8min), then every 6h — the artist-follow cadence
# (was `every 6h` on the gateway). Persistent catches up a tick missed across a box reboot.
# Idempotent by construction (followed_at IS NULL), so a deferred tick just follows the next
# artists on the following cycle.
OnBootSec=8min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
