# docs/agents/hermes/healthcheck-timer/fluncle-healthcheck.service
[Unit]
Description=fluncle-healthcheck — rave-02 host trigger for the Hermes /status prober (decoupled from the busy gateway it monitors)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/healthcheck-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service
# If this unit dies before writing its own /status marker (OOM, missing binary, a crash
# before the sweep catch runs), the healthcheck prober never sees a failure. This catch-all
# posts a minimal Discord alert; %n is this unit's full name.
OnFailure=fluncle-sweep-failure@%n.service

[Service]
Type=oneshot
# The host unit runs as root (to drive the Docker daemon); `docker exec -u hermes` runs the
# probe WORK unprivileged inside the container — matching every other fluncle-* sweep timer,
# so no host cron execs as container-root. The probe WORK is BAKED into the image (the .sh/.ts
# pair at /opt/hermes-scripts/, Unit A — auto-updates from main via pin-watch, no
# /opt/data/scripts copy); this host unit only TRIGGERS it, so the prober no longer rides — and
# is no longer starved by — the Hermes cron gateway it is meant to monitor. (The prober
# self-emits its own /status row, so — unlike the work sweeps — its .sh writes no cron-output
# marker: a self-read would be circular.)
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/fluncle-healthcheck.sh
# A clean tick runs ~5s. The old 120s ceiling could false-kill the prober when it overlapped a
# pin-watch docker image build (~4 min of heavy CPU contention on this 4-vCPU box) — and a
# killed prober fires an OnFailure Discord alert, i.e. contention would masquerade as an
# outage. 300s clears the build window while staying well under the 10-min cadence and short
# enough that a genuinely hung probe still surfaces promptly.
TimeoutStartSec=300
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
