# docs/agents/hermes/label-lineage-timer/fluncle-label-lineage.timer
[Unit]
Description=fluncle-label-lineage timer — fill one bounded batch of label lineage (founding + parent imprint) every 60m on the rave-02 host (not the Hermes gateway)

[Timer]
# First tick shortly after boot (staggered at 13min, off the other sweeps' minutes), then every
# 60m. The `labels` worklist is the findings + crawler-minted labels, and each label is at most two
# MusicBrainz calls at 1 req/s (the CLI circuit-breaks on a throttle), so there is no prize for
# hurrying: an hourly cadence drains the worklist with room to spare. A tick that finds every label
# already walked/attempted is a cheap no-op.
#
# The `labels` worklist is the durable state (`lineage_state` + the `lineage_attempted_at` stamp),
# so a tick missed across a reboot costs nothing — Persistent catches it up and the next tick
# resumes from the top; the state re-skips the labels already walked/attempted.
OnBootSec=13min
OnUnitActiveSec=60min
Persistent=true

[Install]
WantedBy=timers.target
