# docs/agents/hermes/artist-sweep-timer/fluncle-artist-sweep.service
[Unit]
Description=fluncle-artist-sweep — rave-02 host trigger for the artist-resolution sweep (migrated off the Hermes gateway cron runner)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/artist-sweep-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service
# If this unit dies before writing its own /status marker (OOM, missing binary, a crash
# before the sweep catch runs), the healthcheck prober never sees a failure. This catch-all
# posts a minimal Discord alert; %n is this unit's full name.
OnFailure=fluncle-sweep-failure@%n.service

[Service]
Type=oneshot
# Root drives the Docker daemon; `docker exec -u hermes` runs the WORK as the unprivileged
# hermes user. The .sh/.ts pair is BAKED at /opt/hermes-scripts/ (Unit A — auto-updates from
# main via pin-watch); this host unit only TRIGGERS it, and the sweep self-writes its /status
# marker (cron-output.sh) so the prober stays honest.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/artist-sweep.sh
# A pure Worker-trigger (agent-tier resolve_artist: the MusicBrainz url-rel walk + Firecrawl
# gap-fill + YouTube channel-ID resolution); one bounded batch per tick. Zero box tokens. A
# healthy tick is ~55s. The real fix for the chronic timeouts is upstream — every external
# fetch on the resolve path now carries a deadline (the shared MB client, the Firecrawl
# gap-fill, the YouTube lookup), and the box driver bounds each `fluncle` CLI call — so one
# hung item fails that item instead of blocking the whole batch. This 900s ceiling is only
# the backstop: 3× a healthy tick plus the ~4-min pin-watch-build contention window, and far
# under the 60-min cadence.
TimeoutStartSec=900
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
