# docs/agents/hermes/audit-timer/fluncle-audit.timer
[Unit]
Description=fluncle-audit timer — the nightly codebase-audit auditor at 01:00 Amsterdam on the rave-02 host (not the Hermes gateway)

[Timer]
# 01:00 Europe/Amsterdam via the explicit OnCalendar TZ suffix — the host clock is UTC (only the
# Hermes container is TZ-pinned), so we name the zone here; DST-safe, needs systemd >=252. One
# domain per night on a 7-day rotation (audit/rotation.ts). The 5am reviewer (fluncle-audit-review)
# reviews + merges the PR this opens. Persistent catches up a tick missed across a box reboot;
# a late audit is harmless (it just opens tonight's PR whenever it runs).
OnCalendar=*-*-* 01:00:00 Europe/Amsterdam
Persistent=true

[Install]
WantedBy=timers.target
