# docs/agents/hermes/backup-timer/fluncle-backup.timer
[Unit]
Description=fluncle-backup timer — trigger the database-backup sweep nightly at 03:00 Amsterdam on the rave-02 host (not the Hermes gateway)

[Timer]
# 03:00 Europe/Amsterdam via the explicit OnCalendar TZ suffix (the host clock is UTC) — a fixed,
# predictable deep-off-peak nightly slot rather than a rolling 24h interval that drifts across the
# day. It sits before the 05:00 audit-reviewer (which may merge a migration), so the snapshot
# predates it. Persistent catches up a night the box slept through on the next wake (a durable
# off-site backup must not silently skip a day); the prune keeps 30 daily / 12 monthly, so a
# caught-up extra run is harmless.
OnCalendar=*-*-* 03:00:00 Europe/Amsterdam
Persistent=true

[Install]
WantedBy=timers.target
