# docs/agents/hermes/backup-timer/fluncle-backup.timer
[Unit]
Description=fluncle-backup timer — trigger the database-backup sweep every 24h on the rave-02 host (not the Hermes gateway)

[Timer]
# First tick shortly after boot (staggered at 11min), then every 24h — the backup cadence (was
# `every 24h` on the gateway). Persistent catches up a day the box slept through on the next
# wake (a durable off-site backup must not silently skip a day). The prune keeps 30 daily / 12
# monthly, so a caught-up extra run is harmless.
OnBootSec=11min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
