# docs/agents/hermes/frontier-refresh-timer/fluncle-frontier-refresh.timer
[Unit]
Description=fluncle-frontier-refresh timer — one tick of the PACED Frontier-playlist drain, every 15 min on the rave-02 host (not the Hermes gateway)

[Timer]
# The paced drain ticks every 15 minutes. The Frontier still promises "Refreshed weekly", but
# the refresh is now SPREAD across the day instead of bursting every playlist at 07:00: each
# tick processes only a small batch of users whose per-user cursor is older than ~6 days, so the
# whole crew refreshes ~weekly without the one overloaded pass that 429'd Spotify's shared
# per-app budget. A frequent tick is the shape a paced drain needs — the batch size + the ~6-day
# due-gate live in the code (FRONTIER_REFRESH_BATCH / FRONTIER_REFRESH_MIN_AGE_MS), not here.
# Persistent catches up a tick the box slept through on the next boot — harmless, because the
# drain is idempotent (durable per-user cursors), so a late run is just a normal run.
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
