# docs/agents/hermes/frontier-refresh-timer/fluncle-frontier-refresh.service
[Unit]
Description=fluncle-frontier-refresh — rave-02 host trigger for the weekly Frontier playlist refresh (E2, the public recommendation machine)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/frontier-refresh-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
# Root drives the Docker daemon; `docker exec -u hermes` runs the WORK as the unprivileged
# hermes user. The .sh/.ts pair is BAKED at /opt/hermes-scripts/ (auto-updates from main via
# pin-watch); this host unit only TRIGGERS it, and the sweep self-writes its /status marker
# (cron-output.sh). NO new secret: `refresh_frontier_playlists` is AGENT tier, so the box's
# existing agent-scoped token drives it.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/frontier-refresh-sweep.sh
# A PURE trigger — zero LLM tokens. The work is the mirror + the Spotify item replaces inside
# the Worker; one `fluncle admin frontier refresh` walks up to its limit and stops. Give the
# per-user Spotify round-trips headroom as the crew grows.
TimeoutStartSec=900
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
