# docs/agents/hermes/enrich-timer/fluncle-enrich.timer
[Unit]
Description=fluncle-enrich timer — trigger the on-box enrichment sweep every 5m on the rave-02 host (not the Hermes gateway)

[Timer]
# First tick shortly after boot (staggered at 1min so the simultaneous 5-min sweeps don't
# pile), then every 5m — the enrich cadence (was `every 5m` on the gateway). Persistent
# catches up a tick missed across a box reboot on the next wake. The queue is idempotent
# (status=queue), so a deferred tick just enriches on the next cycle.
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
