# docs/agents/hermes/enrich-timer/fluncle-enrich.service
[Unit]
Description=fluncle-enrich — rave-02 host trigger for the on-box enrichment sweep (migrated off the Hermes gateway cron runner)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/enrich-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service
# If this unit dies before writing its own /status marker (OOM, missing binary, a crash
# before the sweep catch runs), the healthcheck prober never sees a failure. This catch-all
# posts a minimal Discord alert; %n is this unit's full name.
OnFailure=fluncle-sweep-failure@%n.service

[Service]
Type=oneshot
# The host unit runs as root (to drive the Docker daemon); `docker exec -u hermes` runs the
# in-container enrichment WORK as the UNPRIVILEGED hermes user (hermes owns the 0600 secrets
# file the .sh sources). The sweep trio is BAKED into the image at /opt/hermes-scripts/ and
# the enrichment DSP skill at /opt/hermes-skills/ (Unit A — both auto-update from main via
# pin-watch, no /opt/data copy); this host unit only TRIGGERS it, and the sweep self-writes
# its /status marker (cron-output.sh) so the fluncle-healthcheck prober stays honest. Parallel
# host timers replace the one serial gateway runner (~300s global budget); a staggered
# OnBootSec spreads the simultaneous 5-min ticks so the latency-sensitive sweeps never pile.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/enrich-sweep.sh
# Pure on-box DSP over a bounded batch (cap 4); a clean tick is seconds, a full batch well
# under two minutes. This is the most CPU-bound sweep, so a pin-watch docker image build (~4
# min of heavy CPU contention on this 4-vCPU box) overlapping a full batch stretches it the
# most — the old 180s ceiling false-killed such a tick. 480s clears that. It exceeds the 5-min
# cadence: systemd never starts a new tick while one is active, so a contention-stretched batch
# just skips the next trigger and self-heals — it is a backstop, never the normal path.
TimeoutStartSec=480
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
