# docs/agents/hermes/secrets/fluncle-secrets-sync.service
[Unit]
Description=Fluncle secrets sync (op inject from 1Password into the box env files)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fluncle-secrets-sync.sh
