# docs/agents/hermes/backfill-timer/fluncle-backfill.service
[Unit]
Description=fluncle-backfill — rave-02 host trigger for the catalogue-backfill sweep (migrated off the Hermes gateway cron runner)
Documentation=https://github.com/mauricekleine/fluncle/blob/main/docs/agents/hermes/backfill-timer/README.md
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
# Root drives the Docker daemon; `docker exec -u hermes` runs the WORK as the unprivileged
# hermes user. The .sh/.ts pair is BAKED at /opt/hermes-scripts/ (Unit A — auto-updates from
# main via pin-watch); this host unit only TRIGGERS it, and the sweep self-writes its /status
# marker (cron-output.sh) so the prober stays honest.
ExecStart=/usr/bin/docker exec -u hermes -e HOME=/opt/data/home hermes bash /opt/hermes-scripts/backfill-sweep.sh
# A pure Worker-trigger (Discogs id + Last.fm love), one paced batch of each source per tick
# (--limit 6). The Worker carries the Retry-After backoff, so the box stays dumb + fast.
TimeoutStartSec=180
# Modest hardening for a root docker-driving oneshot.
NoNewPrivileges=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
